# check_dramatiq.service
[Unit]
Description=Check if dramatiq workers are working

[Service]
User=spekit
Group=spekit
ExecStart=/home/spekit/check_dramatiq.sh

